# btrfs-btrbk-systemd.timer
[Unit]
Description=btrfs-btrbk hourly backup

[Timer]
OnCalendar=hourly
AccuracySec=10min
Persistent=true

[Install]
WantedBy=timers.target
